# scripts/com.mytrainingos.dailysync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mytrainingos.dailysync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/marco/.gemini/antigravity/scratch/mytrainingos/scripts/daily_sync.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/marco/.gemini/antigravity/scratch/mytrainingos/logs/launchd_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/marco/.gemini/antigravity/scratch/mytrainingos/logs/launchd_stderr.log</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
